# srcs/requirements/bonus/mariadb_backup/crontab
*/1 * * * * root /backup.sh >> /var/log/backup.log
